# config/schedule.cron
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

0 6 * * * /app/scripts/run-checkin.sh
6 11 * * * /app/scripts/run-checkin.sh
12 16 * * * /app/scripts/run-checkin.sh
18 21 * * * /app/scripts/run-checkin.sh
